# provision_files/supervisord/websockets.conf
[program:websockets]
command=node /var/www/html/server/run.js
autostart=true
autorestart=true
stderr_logfile=/var/log/websockets.error.log
stdout_logfile=/var/log/websockets.output.log
